# debian/here-server.service
[Unit]
Description=Here Server - Location API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/here-server
Restart=always
RestartSec=5
EnvironmentFile=/etc/here-server/env
DynamicUser=no
User=here-server
Group=here-server
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
