# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    if [ ! -f /etc/here-server/env ]; then
        mkdir -p /etc/here-server
        cat > /etc/here-server/env <<'EOF'
# Here Server - 环境配置
PORT=9001
DATA_DIR=/var/lib/here-server
MAX_HOURS=24
EOF
        chmod 600 /etc/here-server/env
    fi
    mkdir -p /var/lib/here-server
    chmod 700 /var/lib/here-server

    systemctl daemon-reload
    systemctl enable here-server || true
    systemctl start here-server || true
fi
