# debian/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl stop i-am-here || true
    systemctl disable i-am-here || true
fi
