# debian/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl stop here-server || true
    systemctl disable here-server || true
fi
